# wttr.service
[Unit]
Description=Simple weather service

[Service]

ExecStart=/home/michael/bin/wttr.sh

[Install]
WantedBy=multi-user.target
