# wttr.timer
[Unit]
Description=Restarts wttr.service hourly

[Timer]
Unit=wttr.service
OnCalendar=*:*:0/10
Persistent=true

[Install]
WantedBy=timers.target
